# site/profile/files/puppet/puppet-master.timer
[Unit]
Description=Run Puppet master every 30 minutes

[Timer]
OnActiveSec=30min

[Install]
WantedBy=timers.target
