# old/profile/files/bootstrap/node/puppet-master.timer
[Unit]
Description=Run Puppet master every 30 minutes

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w 

[Install]
WantedBy=timers.target
